# php-apache2/services.d/apache2/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/apache2ctl -D FOREGROUND
